# gemini_api.service
[Unit]
Description=Gemini CLI API Server
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/gemini-cli-api
ExecStart=/opt/gemini-cli-api/venv/bin/python /opt/gemini-cli-api/main.py
EnvironmentFile=/opt/gemini-cli-api/.env
Restart=always
RestartSec=5
User=root
Group=root

[Install]
WantedBy=multi-user.target
